# sdc-ntp-clock/files/etc/rc.local
/etc/init.d/sysntpd disable

snmpd -c /sysconf/net-snmp/snmp/snmpd.conf

/SDC/Clock &

exit 0
