# sdc-ntp-clock/files/etc/hotplug.d/iface/77-management
ifconfig eth0:0 10.10.6.10 netmask 255.255.255.0 broadcast 10.10.6.255 up
